# server/supervisor.ini
[program:xl-server]
command=gunicorn -b 0.0.0.0:8108 --chdir /root/server/main main:app
directory=/root/server
autostart=true
autorestart=true
priority=5
killasgroup=true
stopasgroup=true
